# runner/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# Virtual display for the visible browser.
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x720x24 -ac +extension GLX +render -noreset
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# Lightweight window manager so Chromium windows are managed/visible.
[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# Expose the X display over VNC on localhost only (Node gates external access).
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -localhost -quiet
autorestart=true
priority=30
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# websockify serves the noVNC web client and bridges WS -> VNC (localhost:6080).
[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 127.0.0.1:5900
autorestart=true
priority=40
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# The public Node server: /health, /run-tool, /vnc proxy. Binds Render's $PORT.
[program:server]
command=/usr/local/bin/node /app/server.js
directory=/app
environment=DISPLAY=":99"
autorestart=true
priority=50
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
